# scripts/altar-trader-hub.service
[Unit]
Description=Altar Trader Hub Backend API
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=samuelevalente
Group=samuelevalente
WorkingDirectory=/home/samuelevalente/altar-trader-hub-be
Environment="PATH=/home/samuelevalente/altar-trader-hub-be/venv/bin"
Environment="PYTHONUNBUFFERED=1"

# Carica le variabili d'ambiente dal file .env
EnvironmentFile=-/home/samuelevalente/altar-trader-hub-be/.env

# Comando per avviare l'applicazione
ExecStart=/home/samuelevalente/altar-trader-hub-be/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8001 --workers 4

# Riavvia il servizio in caso di errore
Restart=always
RestartSec=10

# Timeout
TimeoutStopSec=30

# Log
StandardOutput=journal
StandardError=journal
SyslogIdentifier=altar-trader-hub

[Install]
WantedBy=multi-user.target

